# backend-rust-microservices/supervisord-all.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Worker starts first
[program:worker]
command=/usr/local/bin/worker
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/worker.err.log
stdout_logfile=/var/log/supervisor/worker.out.log
environment=
    RUST_LOG="%(ENV_RUST_LOG)s",
    MONGO_URI="%(ENV_MONGO_URI)s",
    MONGO_DB_NAME="%(ENV_MONGO_DB_NAME)s",
    REQ_BIND_ADDR="tcp://127.0.0.1:5555",
    PULL_BIND_ADDR="tcp://127.0.0.1:5556",
    PUB_BIND_ADDR="tcp://127.0.0.1:5557"
priority=1
startsecs=5

# Gateway starts after worker
[program:gateway]
command=/usr/local/bin/gateway
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/gateway.err.log
stdout_logfile=/var/log/supervisor/gateway.out.log
environment=
    RUST_LOG="%(ENV_RUST_LOG)s",
    PORT="3000",
    MONGO_URI="%(ENV_MONGO_URI)s",
    MONGO_DB_NAME="%(ENV_MONGO_DB_NAME)s",
    WORKER_REQ_ENDPOINT="tcp://127.0.0.1:5555",
    WORKER_PUSH_ENDPOINT="tcp://127.0.0.1:5556",
    WORKER_SUB_ENDPOINT="tcp://127.0.0.1:5557"
priority=2
depends_on=worker
startsecs=3
